# atg.path
[Unit]
Description=AppImage to Gnome file watcher

[Path]
PathChanged=/home/%u/apps/
Unit=atg.service

[Install]
WantedBy=atg.service